# script/schedule.cron
PATH=/usr/local/bin:/usr/bin:/bin

@hourly  cd search_pocket && script/bootstrap.rb -c config/config.yml -s config/sphinx.conf
